# BR-Remote-Ops/schtasks/sample_task.xml
<?xml version="1.0" ?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
    <Triggers>
        <LogonTrigger>
            <Enabled>true</Enabled>
            <UserId>darkvortex.corp\Administrator</UserId>
        </LogonTrigger>
    </Triggers>
    <Principals>
        <Principal>
            <UserId>darkvortex.corp\Administrator</UserId>
        </Principal>
    </Principals>
    <Settings>
        <AllowStartOnDemand>true</AllowStartOnDemand>
        <Enabled>true</Enabled>
    </Settings>
    <Actions>
        <Exec>
            <Command>notepad.exe</Command>
        </Exec>
    </Actions>
</Task>
